# services.d/freshclam/run
#!/bin/sh
exec freshclam --daemon --foreground=true
